# src/etc/cron.d/backup
0 3 * * * root (source /srv/charlotte-website/.env && /usr/libexec/charlotte-website/charlotte-website backup)