# src/etc/cron.d/renew_certs
0 */12 * * * root certbot renew --post-hook "systemctl reload nginx"